# orchestrator-salt/src/main/resources/salt-common/salt/monitoring/systemd/cdp-node-exporter.service
[Unit]
Description=CDP Node Exporter
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/node_exporter/node_exporter --web.config=/opt/node_exporter/node_exporter-web-config.yml

Restart=always

[Install]
WantedBy=multi-user.target